# ansible/roles/app/templates/macrocoats-api.service.j2
[Unit]
Description=Macro Coats API
After=network.target postgresql.service

[Service]
Type=simple
User={{ app_user }}
WorkingDirectory={{ app_dir }}
EnvironmentFile={{ app_dir }}/.env
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
